# infra/quadlet/tradebot-ai-brain.container
[Unit]
Description=Tradebot AI Brain (strategies + intent generator + LLM)
After=tradebot-trading-core.container
Requires=tradebot-trading-core.container

[Container]
Image=localhost/ai-brain:0.1
ContainerName=ai-brain

Network=tradebot.network
NetworkAlias=ai-brain

Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=NODE_ENV=production
Environment=TZ=Asia/Seoul
Environment=PORT=3001

# Trading Core URL (internal)
Environment=TRADING_CORE_URL=http://trading-core:8080

# Strategy Configuration
Environment=STRATEGY_ID=mom_v1

# Ollama Configuration (local LLM - optional)
Environment=OLLAMA_URL=http://localhost:11434
Environment=OLLAMA_MODEL=llama2

[Service]
Restart=always
RestartSec=5
TimeoutStartSec=300

[Install]
WantedBy=default.target
