# infra/quadlet/tradebot-postgres.container
[Unit]
Description=Tradebot Postgres
After=tradebot.network
Requires=tradebot.network

[Container]
Image=docker.io/library/postgres:16
ContainerName=tradebot-postgres

Network=tradebot.network
NetworkAlias=postgres

Volume=pgdata.volume:/var/lib/postgresql/data

Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=POSTGRES_USER=tradebot
Environment=POSTGRES_PASSWORD=CHANGE_ME_STRONG_PASSWORD
Environment=POSTGRES_DB=tradebot
Environment=TZ=Asia/Seoul

HealthCmd=pg_isready -U tradebot -d tradebot
HealthInterval=10s
HealthTimeout=3s
HealthRetries=10

[Service]
Restart=always
RestartSec=3
TimeoutStartSec=300

[Install]
WantedBy=default.target
